# .build/supervisord/04-cron.conf
[program:cron]
command=cron -f
autostart=true
autorestart=true
priority=99
stdout_events_enabled=true
stderr_events_enabled=true
