# .build/supervisord/06-startup.conf
[program:startup-script]
user=nobody
directory=%(ENV_HOME_SITE)s
command=
autostart=false
autorestart=false
startretries = 0
priority=25
stdout_events_enabled=true
stderr_events_enabled=true
